# dot_config/sv/swaybg/run
#!/bin/sh
test -f /run/user/1000/wayland-*.lock >/dev/null || exit 1

#exec swaybg -i ~/.config/river/wall.png
wall=$(fd . ~/.dotfiles/wallpapers | shuf -n 1)
exec swaybg -i $wall
