# sv/.config/sv/easyeffects/run
#!/bin/sh
exec easyeffects --gapplication-service
